# etc/init.d/mba62
#!/lib/init/init-d-script
# kFreeBSD do not accept scripts as interpreters, using #!/bin/sh and sourcing.
#if [ true != "$INIT_D_SCRIPT_SOURCED" ] ; then
#    set "$0" "$@"; INIT_D_SCRIPT_SOURCED=true . /lib/init/init-d-script
#fi
### BEGIN INIT INFO
# Provides:          mba62
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: MacBookAir6,2 helpers
# Description:       This script runs various helpers for MacBookAir6,2.
### END INIT INFO

# Author: Ville Korhonen <ville@xd.fi>

NAME="mba62"
DESC="MacBookAir6,2 helpers"
DAEMON=none
PIDFILE=none

#

change_device_permissions() {
	DEVICES="$@"
	for device in ${DEVICES}
	do
		echo "Setting device ${device}..."
		ctrl_file="${device}/brightness"
		if [ -e "${ctrl_file}" ]
		then
			chgrp ${MBA62_BACKLIGHT_GROUP} ${ctrl_file}
			chmod g+w ${ctrl_file}
		fi
	done
}

init_brightness() {
	/usr/lib/mba62/init-brightness.sh
}

#
[ -r /etc/default/mba62 ] || exit 0
. /etc/default/mba62

change_device_permissions "${MBA62_ALL_DEVICES}"
init_brightness

exit 0